# deploy/airview-aihot.service
[Unit]
Description=Solution AI Mate site for aihot.muchu.cloud
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/aihot.muchu.cloud/current
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=18173
EnvironmentFile=-/etc/airview-aihot.env
ExecStart=/usr/bin/node /opt/aihot.muchu.cloud/current/server.js
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/aihot.muchu.cloud

[Install]
WantedBy=multi-user.target
